# data/gigamate.service
[Unit]
Description=GigaMate — Gigabyte laptop management
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=exec
ExecStart=%h/.local/bin/gigamate-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
